# bin/rfsom-box-gui.service.in
[Unit]
Description=RFSOM-BOX GUI Application
After=multi-user.target

[Service]
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/analog/.Xauthority
ExecStart=$$INSTALL_LOCATION/bin/rfsom-box-gui-start.sh
Restart=always
RestartSec=10s
KillMode=process
TimeoutSec=infinity

[Install]
WantedBy=graphical.target

